# templates/galaxy-backup.timer.j2
[Unit]
Description=Run galaxy backup periodically

[Timer]
OnCalendar={{ galaxy_backup_calendar }}
Persistent=true
Unit=galaxy-backup.service

[Install]
WantedBy=timers.target